# autoswarm/templates/supervisor.conf.j2
[program:autoswarm]
command=/opt/autoswarm/autoswarm-update-and-start.sh {{swarmrole}}
user={{ansible_user}}
autostart=true
autorestart=true
stopasgroup=true

