# datapool/templates/datapool.mount.j2
[Unit]
Description=Mount the {{ src }} to {{ dest }}

[Mount]
What={{ src }}
Where={{ dest }}
Type=ext4
Options=defaults

[Install]
WantedBy=multi-user.target
