# docker-tcp-switchboard/templates/supervisor.conf.j2
[program:docker-tcp-switchboard]
command=/usr/local/bin/docker-tcp-switchboard.py /etc/docker-tcp-switchboard.conf {{ listen_address }}
user={{ansible_user}}
autostart=true
autorestart=true

